# local_testing/helloworld/infra/deploy_files/local/etc/systemd/system/helloworld.service
[Unit]
Description=Helloworld Gunicorn Django App
After=network.target

[Service]
Type=notify
User=app
Group=app
RuntimeDirectory=helloworld
WorkingDirectory=/home/app/apps/helloworld/current
ExecStart=/home/app/apps/helloworld/current/.venv/bin/gunicorn \
    --workers 2 \
    --bind unix:/home/app/apps/helloworld/current/helloworld.sock \
    --access-logfile - \
    --error-logfile - \
    --timeout 30 \
    helloworld.wsgi:application
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
